# .config/VSCodium/User/History/-560f93d9/FK2d.timer
[Unit]
Description=Weekly script to maintain systems
# Requires=qbit-tracker-add.service

[Timer]
OnCalendar=Weekly
RandomizedDelaySec=120
Persistent=true

[Install]
WantedBy=timers.target
